# templates/benchmark-metrics-exporter.service.j2
[Unit]
Description=Benchmark Metrics Exporter
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Restart=no

Environment=ETCD_BENCHMARK_JSONL_PATH={{ jsonl_path }}
Environment=ETCD_BENCHMARK_TEST_ID={{ test_id }}
Environment=ETCD_BENCHMARK_METRICS_PORT={{ metrics_port }}
Environment=ETCD_BENCHMARK_REPORT_INTERVAL={{ report_interval }}

ExecStart={{ python_bin }} {{ runner_path }}

StandardOutput=journal
StandardError=journal